# scripts/blocklist/blocklist-restore.service
[Unit]
Description=Restore blocklist ipsets at boot (before UFW)
DefaultDependencies=no
# Ensure /var/lib/blocklist is mounted before we try to read save files
After=local-fs.target
RequiresMountsFor=/var/lib/blocklist
Before=network-pre.target
Before=ufw.service
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/restore-blocklist.sh

[Install]
WantedBy=multi-user.target
