# scripts/blocklist/blocklist.service
[Unit]
Description=Update abuse IP blocklist from fail2ban logs
# Wants= (not Requires=) — run even if fail2ban is temporarily down so the
# updater can still refresh metrics and process rotated logs.
Wants=fail2ban.service
After=fail2ban.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-blocklist.sh
# Service hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# Tight write paths — only what the updater actually needs
ReadWritePaths=/var/lib/blocklist /var/lib/node_exporter/textfile_collector
ReadWritePaths=/var/log/blocklist-updates.log
